# ubuntu/data/post-install.service
[Unit]
Description=My Script
After=network.target
[Service]
ExecStart=/usr/share/drk/post-install.sh
[Install]
WantedBy=default.target